# ARCHIVE/taknet-ps-v2.8-network-features/systemd/adsb-web.service
[Unit]
Description=TAKNET-PS-ADSB Web Interface
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/adsb/web
ExecStart=/usr/bin/python3 /opt/adsb/web/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
